# supplement/inittab
::sysinit:/etc/init.d/rc
::askfirst:/bin/sh
::restart:/sbin/init
::ctrlaltdel:/sbin/reboot -f
::shutdown:/bin/umount -a -r
